# OllamaMacOSApp/com.ollamamacos.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ollamamacos.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/OllamaMacOSApp.app/Contents/MacOS/ollama-daemon</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/ollama-daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/ollama-daemon-error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Applications/OllamaMacOSApp.app/Contents/Resources</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>OLLAMA_HOST</key>
        <string>127.0.0.1:11434</string>
        <key>OLLAMA_MODELS</key>
        <string>/Applications/OllamaMacOSApp.app/Contents/Resources/models</string>
    </dict>
</dict>
</plist>
